# containers/astra/atse-full/config-mgmt/etc/profile.d/allinea.csh
#!/bin/csh

setenv ALLINEA_LICENSE_DIR /opt/arm/forge/licenses
