# containers/astra/atse-full/config-mgmt/etc/profile.d/allinea.sh
#!/bin/sh

export ALLINEA_LICENSE_DIR=/opt/arm/forge/licenses
